# deployment/server_deploy/run-mongo.conf
description "a init script for mongod service"
author "Pepe Lu"

respawn
respawn limit 10 60

start on started sshd
stop on shutdown

script
  exec mongod -f /etc/mongod.conf
end script
